# ramfs/sdk/prefix/etc/inittab
::sysinit:/etc/rc.d/rcS start
::respawn:/bin/sh < /dev/ttyAMA0 2>&1 > /dev/ttyAMA0
::restart:/sbin/init
::ctrlaltdel:/sbin/shutdown -t3 -r now
::shutdown:/etc/rc.d/rcK
